# chronyd-ptp/src/chronyd-restart.service
[Unit]
Description=Restart chronyd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restart-chronyd
